# ghost-ses-adapter.service
[Unit]
Description=Ghost SES Adapter
After=network.target

[Service]
ExecStart=/usr/bin/node /path/to/ghost-ses-adapter.js
WorkingDirectory=/path/to/ghost-ses-adapter
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ghost-ses-adapter
User=your-user
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
